# systemd/msmtpctl.service
[Unit]
Description=Reload msmtp config service

Requires=msmtpctl.timer

[Service]
Type=oneshot
ExecStart=%h/bin/msmtpctl --reload
